# fnkeys/zenbook-duo-systools-fnkeys.service
[Unit]
Description=Zenbook Duo Fn-key user-session helper
After=default.target
ConditionPathExists=/usr/bin/zenbook-duo-systools-fnkeys

[Service]
Type=simple
ExecStart=/usr/bin/zenbook-duo-systools-fnkeys
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
