# sysstates/zenbook-duo-systools-user.service
[Unit]
Description=Zenbook Duo user-session display watcher
After=default.target
ConditionPathExists=/usr/bin/zenbook-duo-systools

[Service]
Type=simple
ExecStart=/usr/bin/zenbook-duo-systools watch
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target
