# sysstates/zenbook-duo-systools-user.service
[Unit]
Description=Zenbook Duo user-session display tools
After=default.target
ConditionPathExists=/usr/bin/zenbook-duo-systools

[Service]
Type=oneshot
ExecStart=/usr/bin/zenbook-duo-systools apply

[Install]
WantedBy=default.target
